# CONTAINER/citadel_jugo.container
[Unit]
Description=CITADEL repo app jugo (citadel)

[Container]
ContainerName=citadel-jugo
HostName=citadel-jugo
Image=localhost/citadel-jugo:latest
Network=host
AddCapability=NET_ADMIN
AddCapability=NET_RAW
AddDevice=/dev/net/tun
EnvironmentFile=/home/openclaw/saf/JUGO/CONTAINER/citadel_jugo.env
Volume=citadel-jugo-var-lib-tailscale:/var/lib/tailscale:Z
# PublishPort lines are ignored in host mode

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
